# agent/config/restarter.service
[Unit]
Description=restarts dead containers
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/appmgr repair-app-status