# agent/config/restarter.timer
[Unit]
Description=restarter

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target