# agent/config/restarter.service
[Unit]
Description=restarts dead containers
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash /root/restarter.sh