# build/initialization.service
[Unit]
Description=Boot process for system initialization.
After=network-online.target systemd-time-wait-sync.service

[Service]
Type=oneshot
Restart=on-failure
RestartSec=5s
ExecStart=/bin/sh -c "initialization.sh | nc-broadcast --tee 0.0.0.0:8080"
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
