# build/initialization.service
[Unit]
Description=Boot process for system initialization.
After=network-online.target systemd-time-wait-sync.service

[Service]
Type=oneshot
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/initialization.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
