# container-runtime/container-runtime-failure.service
[Unit]
Description=StartOS Container Runtime Failure Handler

[Service]
Type=oneshot
ExecStart=/usr/bin/start-cli rebuild